# roles/repositories/files/paccache-cleanup.timer
[Unit]
Description=paccache cleanup service timer

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=True

[Install]
WantedBy=timers.target
